# supervisord.conf
[supervisord]
nodaemon=true

[program:monogod]
command=service mongod start
autorestart=true

[program:node]
command=node DormitoryD/index.js
autorestart=true